# supervisord/github_runner.ini
[program:runner]
process_name=runner
directory=/home/centos/actions-runner
command=/home/centos/actions-runner/run.sh
user=centos
numprocs=1
autostart=true
autorestart=true
redirect_stderr=true
stopsignal=TERM
stdout_logfile_maxbytes=5GB
stdout_logfile_backups=300
stdout_logfile=/home/centos/actions-runner/runner_stdout.log
stderr_logfile=/home/centos/actions-runner/runner_stderr.log
